# rust/packaging/claude-usage-indicator-dev.service
# dev systemd --user 单元：开机自启 dev 构建的 Rust 托盘，与 prod 正式版并存（独立 APP_ID）。
# 注意：dev 与 prod 的 release 产物同路径 target/release/cui，谁后编谁覆盖。本单元要跑 dev，
#   必须确保该二进制是用 `cargo build -p cui --release --features dev` 编的（--version 末尾带 -dev）。
#   迭代开发更推荐直接 `cargo run -p cui --features dev`（实时日志、秒重启），无需常驻单元。
# 安装：cp 到 ~/.config/systemd/user/ ，按需改 ExecStart 路径，再
#   systemctl --user daemon-reload && systemctl --user enable --now claude-usage-indicator-dev.service
[Unit]
Description=Claude Usage Indicator (dev)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/claude-usage-indicator/rust/target/release/cui
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
