# packaging/claude-usage-indicator.service
[Unit]
Description=Claude Usage Indicator
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/run.sh
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
